# agent/charms/testflinger-agent-host-charm/templates/testflinger-agent.supervisord.conf.j2
[program:{{ agent_name }}]
redirect_stderr=true
environment=PYTHONIOENCODING=utf-8
user=ubuntu
command={{ virtual_env_path }}/bin/testflinger-agent -c {{ agent_config_path }}/testflinger-agent.conf
