# agent/charms/testflinger-agent-host-charm/tests/unit/test_data/agent_no_port.conf
[program:agent-no-port]
startsecs=0
redirect_stderr=true
user=ubuntu
command=/srv/testflinger-venv/bin/testflinger-agent -c /srv/config/test.conf