# agent/charms/testflinger-agent-host-charm/tests/unit/test_data/agent2.conf
[program:agent2]
startsecs=0
redirect_stderr=true
environment=USER="ubuntu",HOME="/home/ubuntu",PYTHONIOENCODING=utf-8
user=ubuntu
command=/srv/testflinger-venv/bin/testflinger-agent -c /srv/agent-configs/ce-cert/data-tel-l2/agent2/testflinger-agent.conf -p 8001